# examples/openrc/wtmpd
#!/sbin/openrc-run

name="wtmpd"
command="s6-ipcserver"
command_args="/run/utmps/wtmpd-socket utmps-wtmpd"
command_background=yes
command_user=utmp
pidfile="/run/utmps/wtmpd.pid"
start_stop_daemon_args="-d /run/utmps"

start_pre() {
   checkpath -D -d -o utmp:utmp -m 0755 /run/utmps
}
